# workspace/crontab/cash-track
* * * * * finance php /var/www/artisan schedule:run > /var/log/cron-job.log 2>&1